# deploy/pdf-audio-backend.service
[Unit]
Description=Bilingual PDF Audio - FastAPI Backend
After=network.target redis-server.service

[Service]
Type=simple
User=pdf-audio
Group=pdf-audio
WorkingDirectory=/opt/pdf-audio/backend
Environment=APP_ENV=production
Environment=PYTHONPATH=/opt/pdf-audio/backend
Environment=PATH=/opt/pdf-audio/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/pdf-audio/pdf-audio.env
ExecStart=/opt/pdf-audio/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8543
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/pdf-audio /opt/pdf-audio/storage

[Install]
WantedBy=multi-user.target
